# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tarot-backend]
command=/usr/local/bin/tarot-backend
autostart=true
autorestart=true
environment=DEEPSEEK_API_KEY=%(ENV_DEEPSEEK_API_KEY)s,DEEPSEEK_API_URL=%(ENV_DEEPSEEK_API_URL)s,PORT=8081
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
